# support/ithowp_request.service
[Unit]
Description=Request Datalog message from Itho WPU

[Service]
Type=oneshot
User=root
ExecStart=/home/pi/pislave/pimaster 97
